# rootfs/etc/services.d/version/run
#!/usr/bin/with-contenv sh

exec /bin/version -port=80 -hosts=$VERSION_HOSTS
